# overlay-pi/etc/rc.local
#!/bin/bash

dmesg --console-off
exit 0
